# supervisor/webapp.conf
[program:webapp]
command=/usr/bin/gunicorn3 app:app
directory=/home/pi/test-update
autostart=true
autorestart=true
startretries=3
stderr_logfile=/home/pi/test-update/logs/webapp.err.log
stdout_logfile=/home/pi/test-update/logs/webapp.out.log
user=root
